# cron.daily/space_temp_cleanup
#!/bin/sh

find /space/Unsorted/Temp -mtime -14 -exec rm -rf {} \;

